# sysutils/salt/PLIST
@comment $NetBSD: PLIST,v 1.3 2012/12/06 08:01:21 apb Exp $
${PYSITELIB}/salt/__init__.py
${PYSITELIB}/salt/_compat.py
${PYSITELIB}/salt/exceptions.py
${PYSITELIB}/salt/client.py
${PYSITELIB}/salt/config.py
${PYSITELIB}/salt/crypt.py
${PYSITELIB}/salt/fileclient.py
${PYSITELIB}/salt/loader.py
${PYSITELIB}/salt/key.py
${PYSITELIB}/salt/master.py
${PYSITELIB}/salt/log.py
${PYSITELIB}/salt/minion.py
${PYSITELIB}/salt/overstate.py
${PYSITELIB}/salt/payload.py
${PYSITELIB}/salt/runner.py
${PYSITELIB}/salt/scripts.py
${PYSITELIB}/salt/state.py
${PYSITELIB}/salt/template.py
${PYSITELIB}/salt/version.py
${PYSITELIB}/salt/cli/__init__.py
${PYSITELIB}/salt/cli/batch.py
${PYSITELIB}/salt/cli/caller.py
${PYSITELIB}/salt/cli/cp.py
${PYSITELIB}/salt/ext/__init__.py
${PYSITELIB}/salt/auth/__init__.py
${PYSITELIB}/salt/auth/keystone.py
${PYSITELIB}/salt/auth/ldap.py
${PYSITELIB}/salt/auth/pam.py
${PYSITELIB}/salt/wheel/file_roots.py
${PYSITELIB}/salt/wheel/__init__.py
${PYSITELIB}/salt/wheel/config.py
${PYSITELIB}/salt/wheel/pillar_roots.py
${PYSITELIB}/salt/wheel/key.py
${PYSITELIB}/salt/tops/__init__.py
${PYSITELIB}/salt/tops/cobbler.py
${PYSITELIB}/salt/tops/ext_nodes.py
${PYSITELIB}/salt/tops/mongo.py
${PYSITELIB}/salt/grains/__init__.py
${PYSITELIB}/salt/grains/core.py
${PYSITELIB}/salt/grains/extra.py
${PYSITELIB}/salt/grains/opts.py
${PYSITELIB}/salt/modules/alternatives.py
${PYSITELIB}/salt/modules/__init__.py
${PYSITELIB}/salt/modules/aliases.py
${PYSITELIB}/salt/modules/augeas_cfg.py
${PYSITELIB}/salt/modules/apache.py
${PYSITELIB}/salt/modules/apt.py
${PYSITELIB}/salt/modules/archive.py
${PYSITELIB}/salt/modules/at.py
${PYSITELIB}/salt/modules/butterkvm.py
${PYSITELIB}/salt/modules/bluez.py
${PYSITELIB}/salt/modules/brew.py
${PYSITELIB}/salt/modules/daemontools.py
${PYSITELIB}/salt/modules/cassandra.py
${PYSITELIB}/salt/modules/cmdmod.py
${PYSITELIB}/salt/modules/config.py
${PYSITELIB}/salt/modules/cp.py
${PYSITELIB}/salt/modules/cron.py
${PYSITELIB}/salt/modules/debconfmod.py
${PYSITELIB}/salt/modules/data.py
${PYSITELIB}/salt/modules/disk.py
${PYSITELIB}/salt/modules/freebsd_sysctl.py
${PYSITELIB}/salt/modules/debian_service.py
${PYSITELIB}/salt/modules/hg.py
${PYSITELIB}/salt/modules/dnsmasq.py
${PYSITELIB}/salt/modules/djangomod.py
${PYSITELIB}/salt/modules/dnsutil.py
${PYSITELIB}/salt/modules/dpkg.py
${PYSITELIB}/salt/modules/ebuild.py
${PYSITELIB}/salt/modules/eix.py
${PYSITELIB}/salt/modules/event.py
${PYSITELIB}/salt/modules/extfs.py
${PYSITELIB}/salt/modules/file.py
${PYSITELIB}/salt/modules/freebsdservice.py
${PYSITELIB}/salt/modules/freebsdjail.py
${PYSITELIB}/salt/modules/freebsdkmod.py
${PYSITELIB}/salt/modules/freebsdpkg.py
${PYSITELIB}/salt/modules/glance.py
${PYSITELIB}/salt/modules/gem.py
${PYSITELIB}/salt/modules/git.py
${PYSITELIB}/salt/modules/grains.py
${PYSITELIB}/salt/modules/gentoo_service.py
${PYSITELIB}/salt/modules/gentoolkitmod.py
${PYSITELIB}/salt/modules/grub_legacy.py
${PYSITELIB}/salt/modules/groupadd.py
${PYSITELIB}/salt/modules/iptables.py
${PYSITELIB}/salt/modules/guestfs.py
${PYSITELIB}/salt/modules/hosts.py
${PYSITELIB}/salt/modules/img.py
${PYSITELIB}/salt/modules/keyboard.py
${PYSITELIB}/salt/modules/key.py
${PYSITELIB}/salt/modules/linux_sysctl.py
${PYSITELIB}/salt/modules/keystone.py
${PYSITELIB}/salt/modules/kmod.py
${PYSITELIB}/salt/modules/kvm_hyper.py
${PYSITELIB}/salt/modules/launchctl.py
${PYSITELIB}/salt/modules/layman.py
${PYSITELIB}/salt/modules/ldapmod.py
${PYSITELIB}/salt/modules/linux_acl.py
${PYSITELIB}/salt/modules/linux_lvm.py
${PYSITELIB}/salt/modules/openbsdpkg.py
${PYSITELIB}/salt/modules/locale.py
${PYSITELIB}/salt/modules/locate.py
${PYSITELIB}/salt/modules/logrotate.py
${PYSITELIB}/salt/modules/makeconf.py
${PYSITELIB}/salt/modules/mdadm.py
${PYSITELIB}/salt/modules/mongodb.py
${PYSITELIB}/salt/modules/monit.py
${PYSITELIB}/salt/modules/moosefs.py
${PYSITELIB}/salt/modules/mount.py
${PYSITELIB}/salt/modules/munin.py
${PYSITELIB}/salt/modules/mysql.py
${PYSITELIB}/salt/modules/network.py
${PYSITELIB}/salt/modules/nfs3.py
${PYSITELIB}/salt/modules/nginx.py
${PYSITELIB}/salt/modules/nova.py
${PYSITELIB}/salt/modules/npm.py
${PYSITELIB}/salt/modules/nzbget.py
${PYSITELIB}/salt/modules/rpm.py
${PYSITELIB}/salt/modules/openbsdservice.py
${PYSITELIB}/salt/modules/osxdesktop.py
${PYSITELIB}/salt/modules/pacman.py
${PYSITELIB}/salt/modules/pam.py
${PYSITELIB}/salt/modules/parted.py
${PYSITELIB}/salt/modules/pecl.py
${PYSITELIB}/salt/modules/pillar.py
${PYSITELIB}/salt/modules/pip.py
${PYSITELIB}/salt/modules/pkg_resource.py
${PYSITELIB}/salt/modules/pkgng.py
${PYSITELIB}/salt/modules/pkgutil.py
${PYSITELIB}/salt/modules/postgres.py
${PYSITELIB}/salt/modules/poudriere.py
${PYSITELIB}/salt/modules/ps.py
${PYSITELIB}/salt/modules/publish.py
${PYSITELIB}/salt/modules/puppet.py
${PYSITELIB}/salt/modules/pw_group.py
${PYSITELIB}/salt/modules/pw_user.py
${PYSITELIB}/salt/modules/qemu_img.py
${PYSITELIB}/salt/modules/qemu_nbd.py
${PYSITELIB}/salt/modules/quota.py
${PYSITELIB}/salt/modules/rabbitmq.py
${PYSITELIB}/salt/modules/reg.py
${PYSITELIB}/salt/modules/ret.py
${PYSITELIB}/salt/modules/rh_ip.py
${PYSITELIB}/salt/modules/solaris_shadow.py
${PYSITELIB}/salt/modules/rh_service.py
${PYSITELIB}/salt/modules/rvm.py
${PYSITELIB}/salt/modules/s3.py
${PYSITELIB}/salt/modules/saltutil.py
${PYSITELIB}/salt/modules/selinux.py
${PYSITELIB}/salt/modules/service.py
${PYSITELIB}/salt/modules/shadow.py
${PYSITELIB}/salt/modules/smf.py
${PYSITELIB}/salt/modules/solaris_group.py
${PYSITELIB}/salt/modules/solaris_user.py
${PYSITELIB}/salt/modules/solarispkg.py
${PYSITELIB}/salt/modules/solr.py
${PYSITELIB}/salt/modules/sqlite3.py
${PYSITELIB}/salt/modules/ssh.py
${PYSITELIB}/salt/modules/state.py
${PYSITELIB}/salt/modules/status.py
${PYSITELIB}/salt/modules/supervisord.py
${PYSITELIB}/salt/modules/svn.py
${PYSITELIB}/salt/modules/sysbench.py
${PYSITELIB}/salt/modules/sysmod.py
${PYSITELIB}/salt/modules/system.py
${PYSITELIB}/salt/modules/systemd.py
${PYSITELIB}/salt/modules/test.py
${PYSITELIB}/salt/modules/timezone.py
${PYSITELIB}/salt/modules/tls.py
${PYSITELIB}/salt/modules/tomcat.py
${PYSITELIB}/salt/modules/upstart.py
${PYSITELIB}/salt/modules/useradd.py
${PYSITELIB}/salt/modules/virt.py
${PYSITELIB}/salt/modules/virtualenv.py
${PYSITELIB}/salt/modules/win_disk.py
${PYSITELIB}/salt/modules/win_file.py
${PYSITELIB}/salt/modules/win_groupadd.py
${PYSITELIB}/salt/modules/win_network.py
${PYSITELIB}/salt/modules/win_pkg.py
${PYSITELIB}/salt/modules/win_service.py
${PYSITELIB}/salt/modules/win_shadow.py
${PYSITELIB}/salt/modules/win_status.py
${PYSITELIB}/salt/modules/win_useradd.py
${PYSITELIB}/salt/modules/yumpkg.py
${PYSITELIB}/salt/modules/yumpkg5.py
${PYSITELIB}/salt/modules/zpool.py
${PYSITELIB}/salt/modules/zypper.py
${PYSITELIB}/salt/modules/rh_ip/network.jinja
${PYSITELIB}/salt/modules/rh_ip/conf.jinja
${PYSITELIB}/salt/modules/rh_ip/rh5_eth.jinja
${PYSITELIB}/salt/modules/rh_ip/rh6_eth.jinja
${PYSITELIB}/salt/pillar/pillar_ldap.py
${PYSITELIB}/salt/pillar/__init__.py
${PYSITELIB}/salt/pillar/cmd_json.py
${PYSITELIB}/salt/pillar/cmd_yaml.py
${PYSITELIB}/salt/pillar/cobbler.py
${PYSITELIB}/salt/pillar/hiera.py
${PYSITELIB}/salt/pillar/libvirt.py
${PYSITELIB}/salt/pillar/mongo.py
${PYSITELIB}/salt/pillar/puppet.py
${PYSITELIB}/salt/renderers/__init__.py
${PYSITELIB}/salt/renderers/jinja.py
${PYSITELIB}/salt/renderers/json.py
${PYSITELIB}/salt/renderers/mako.py
${PYSITELIB}/salt/renderers/py.py
${PYSITELIB}/salt/renderers/pydsl.py
${PYSITELIB}/salt/renderers/stateconf.py
${PYSITELIB}/salt/renderers/wempy.py
${PYSITELIB}/salt/renderers/yaml.py
${PYSITELIB}/salt/returners/carbon_return.py
${PYSITELIB}/salt/returners/__init__.py
${PYSITELIB}/salt/returners/mongo_future_return.py
${PYSITELIB}/salt/returners/cassandra_return.py
${PYSITELIB}/salt/returners/local.py
${PYSITELIB}/salt/returners/mongo_return.py
${PYSITELIB}/salt/returners/mysql.py
${PYSITELIB}/salt/returners/postgres.py
${PYSITELIB}/salt/returners/redis_return.py
${PYSITELIB}/salt/returners/sentry_return.py
${PYSITELIB}/salt/returners/smtp_return.py
${PYSITELIB}/salt/returners/syslog_return.py
${PYSITELIB}/salt/runners/fileserver.py
${PYSITELIB}/salt/runners/__init__.py
${PYSITELIB}/salt/runners/cache.py
${PYSITELIB}/salt/runners/doc.py
${PYSITELIB}/salt/runners/launchd.py
${PYSITELIB}/salt/runners/jobs.py
${PYSITELIB}/salt/runners/manage.py
${PYSITELIB}/salt/runners/network.py
${PYSITELIB}/salt/runners/search.py
${PYSITELIB}/salt/runners/state.py
${PYSITELIB}/salt/runners/virt.py
${PYSITELIB}/salt/runners/winrepo.py
${PYSITELIB}/salt/states/alternatives.py
${PYSITELIB}/salt/states/__init__.py
${PYSITELIB}/salt/states/alias.py
${PYSITELIB}/salt/states/debconfmod.py
${PYSITELIB}/salt/states/cmd.py
${PYSITELIB}/salt/states/cron.py
${PYSITELIB}/salt/states/grains.py
${PYSITELIB}/salt/states/file.py
${PYSITELIB}/salt/states/gem.py
${PYSITELIB}/salt/states/git.py
${PYSITELIB}/salt/states/keyboard.py
${PYSITELIB}/salt/states/group.py
${PYSITELIB}/salt/states/hg.py
${PYSITELIB}/salt/states/host.py
${PYSITELIB}/salt/states/layman.py
${PYSITELIB}/salt/states/kmod.py
${PYSITELIB}/salt/states/mongodb_database.py
${PYSITELIB}/salt/states/libvirt.py
${PYSITELIB}/salt/states/locale.py
${PYSITELIB}/salt/states/makeconf.py
${PYSITELIB}/salt/states/module.py
${PYSITELIB}/salt/states/mount.py
${PYSITELIB}/salt/states/mongodb_user.py
${PYSITELIB}/salt/states/pkgrepo.py
${PYSITELIB}/salt/states/npm.py
${PYSITELIB}/salt/states/mysql_database.py
${PYSITELIB}/salt/states/mysql_grants.py
${PYSITELIB}/salt/states/mysql_user.py
${PYSITELIB}/salt/states/network.py
${PYSITELIB}/salt/states/pecl.py
${PYSITELIB}/salt/states/pip.py
${PYSITELIB}/salt/states/pkg.py
${PYSITELIB}/salt/states/pkgng.py
${PYSITELIB}/salt/states/postgres_database.py
${PYSITELIB}/salt/states/postgres_group.py
${PYSITELIB}/salt/states/postgres_user.py
${PYSITELIB}/salt/states/rabbitmq_user.py
${PYSITELIB}/salt/states/rabbitmq_vhost.py
${PYSITELIB}/salt/states/rvm.py
${PYSITELIB}/salt/states/selinux.py
${PYSITELIB}/salt/states/service.py
${PYSITELIB}/salt/states/ssh_auth.py
${PYSITELIB}/salt/states/ssh_known_hosts.py
${PYSITELIB}/salt/states/stateconf.py
${PYSITELIB}/salt/states/supervisord.py
${PYSITELIB}/salt/states/svn.py
${PYSITELIB}/salt/states/sysctl.py
${PYSITELIB}/salt/states/user.py
${PYSITELIB}/salt/states/timezone.py
${PYSITELIB}/salt/states/tomcat.py
${PYSITELIB}/salt/states/virtualenv.py
${PYSITELIB}/salt/fileserver/__init__.py
${PYSITELIB}/salt/fileserver/gitfs.py
${PYSITELIB}/salt/fileserver/roots.py
${PYSITELIB}/salt/search/whoosh_search.py
${PYSITELIB}/salt/search/__init__.py
${PYSITELIB}/salt/output/overstatestage.py
${PYSITELIB}/salt/output/__init__.py
${PYSITELIB}/salt/output/grains.py
${PYSITELIB}/salt/output/highstate.py
${PYSITELIB}/salt/output/json_out.py
${PYSITELIB}/salt/output/key.py
${PYSITELIB}/salt/output/nested.py
${PYSITELIB}/salt/output/no_out.py
${PYSITELIB}/salt/output/no_return.py
${PYSITELIB}/salt/output/pprint_out.py
${PYSITELIB}/salt/output/raw.py
${PYSITELIB}/salt/output/txt.py
${PYSITELIB}/salt/output/virt_query.py
${PYSITELIB}/salt/output/yaml_out.py
${PYSITELIB}/salt/utils/atomicfile.py
${PYSITELIB}/salt/utils/__init__.py
${PYSITELIB}/salt/utils/dictupdate.py
${PYSITELIB}/salt/utils/debug.py
${PYSITELIB}/salt/utils/filebuffer.py
${PYSITELIB}/salt/utils/event.py
${PYSITELIB}/salt/utils/gzip_util.py
${PYSITELIB}/salt/utils/find.py
${PYSITELIB}/salt/utils/migrations.py
${PYSITELIB}/salt/utils/jinja.py
${PYSITELIB}/salt/utils/mako.py
${PYSITELIB}/salt/utils/socket_util.py
${PYSITELIB}/salt/utils/minions.py
${PYSITELIB}/salt/utils/parsers.py
${PYSITELIB}/salt/utils/process.py
${PYSITELIB}/salt/utils/pydsl.py
${PYSITELIB}/salt/utils/s3.py
${PYSITELIB}/salt/utils/schedule.py
${PYSITELIB}/salt/utils/verify.py
${PYSITELIB}/salt/utils/saltminionservice.py
${PYSITELIB}/salt/utils/winservice.py
${PYSITELIB}/salt/utils/templates.py
${PYSITELIB}/salt/utils/winapi.py
${PYSITELIB}/salt/utils/yamlloader.py
${PYSITELIB}/salt/utils/xmlutil.py
${PYSITELIB}/salt/_version.py
bin/salt-master
bin/salt-minion
bin/salt-syndic
bin/salt-key
bin/salt-cp
bin/salt-call
bin/salt-run
bin/salt
share/examples/rc.d/salt_master
share/examples/rc.d/salt_minion
share/examples/salt/master
share/examples/salt/minion
${PKGMANDIR}/man1/salt-call.1
${PKGMANDIR}/man1/salt-cp.1
${PKGMANDIR}/man1/salt-key.1
${PKGMANDIR}/man1/salt-master.1
${PKGMANDIR}/man1/salt-minion.1
${PKGMANDIR}/man1/salt-run.1
${PKGMANDIR}/man1/salt-syndic.1
${PKGMANDIR}/man1/salt.1
${PKGMANDIR}/man7/salt.7

